# dsfd.service
[Unit]
Description=Distributed Service Framework (DSF) Daemon

[Service]
EnvironmentFile=/etc/dsf/dsfd.env

User=dsfd
ExecStart=/usr/bin/dsfd

Restart=on-failure

[Install]
WantedBy=multi-user.target
